# templates/udb3/entry_api/udb3-event-export-workers.target.erb
[Unit]
Description=UiTdatabank event export workers
Wants=<%= Array(1..@event_export_worker_count).map { |num| "udb3-event-export-worker@#{num}" }.join(' ') %>

[Install]
WantedBy=multi-user.target
